# redshift/.config/systemd/user/redshift-gtk.service.d/override.conf
# https://github.com/jonls/redshift/issues/265#issuecomment-338245748
[Install]
WantedBy=graphical-session.target
