# mpd/.config/systemd/user/mpDris2.service.d/override.conf
[Unit]
After=mpd.service
Requires=mpd.service
