# tools/.config/systemd/user/dump-lutris-playtime.timer
[Unit]
Description=Record playtimes in Lutris every day (at 0600)

[Timer]
OnCalendar=06:00
Persistent=true

[Install]
WantedBy=timers.target
