# tools/.config/systemd/user/dump-lutris-playtime.timer
[Unit]
Description=Record playtimes every day at 0600 or after boot

[Timer]
OnCalendar=06:00
Persistent=true

[Install]
WantedBy=default.target
